# service/.s6-svscan/finish
#!/bin/bash
set -euo pipefail
echo "s6: finish" 1>&2
